# .config/systemd/user/mount-usb-monitor.service
[Unit]
Description=USB Monitor Service

[Service]
ExecStart=%h/.local/bin/mount-usb-monitor

