# .config/systemd/user/mount-usb-monitor.service
[Unit]
Description=USB mount Monitor Service
After=local-fs.target

[Service]
ExecStart=%h/.local/bin/mount-usb-monitor
Restart=on-failure

[Install]
WantedBy=default.target
